# tests/fixtures/uac_sample/etc/systemd/system/kworker.service
[Unit]
Description=Kernel Worker

[Service]
ExecStart=/tmp/kworker --config /tmp/kworker.conf

[Install]
WantedBy=multi-user.target

